# pmanager/pmanager_update.timer
[Unit]
Description=PManager databases updates

[Timer]
OnCalendar=*.00/15
Persistent=true
