# pmanager/pmanager_update.timer
[Unit]
Description=PManager databases updates

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
