# initscripts-systemd/chakra-daemons.target
[Unit]
Description=Legacy Chakra Daemons

[Install]
WantedBy=multi-user.target
